# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mongodb]
command=/usr/bin/mongod

[program:uhoh]
command=/usr/bin/nodejs /opt/uhoh/index.js